# etc/cron.daily/quota.disabled
#! /bin/sh

# check if quota package is available
test -x /usr/sbin/warnquota || exit 0

# check if warnquota run is configured
test -f /etc/default/quota || exit 0
. /etc/default/quota

if [ "$run_warnquota" = "true" ]; then
	# check if quotas are enabled
	if grep -q '^[^#]*quota' /etc/fstab; then
        	/usr/sbin/warnquota -ug
	fi	
fi

exit 0
